# scripts/services/llms-api.service.example
[Unit]
Description=LLM Scheduling Management System API
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/llm-scheduling-management-system
EnvironmentFile=/opt/llm-scheduling-management-system/.env
ExecStart=/opt/llm-scheduling-management-system/scripts/services/run-api.sh
Restart=always
RestartSec=5
StandardOutput=append:/opt/llm-scheduling-management-system/logs/api.service.log
StandardError=append:/opt/llm-scheduling-management-system/logs/api.service.err.log

[Install]
WantedBy=multi-user.target
